# deploy/predict-market-dashboard.service
[Unit]
Description=Predict Market Bot — Human-readable Dashboard
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=ssm-user
WorkingDirectory=/home/ssm-user/.openclaw/workspaces/predict-market-bot
EnvironmentFile=/home/ssm-user/.openclaw/workspaces/predict-market-bot/.env
ExecStart=/home/ssm-user/.openclaw/workspaces/predict-market-bot/.venv/bin/python dashboard_server.py
Restart=on-failure
RestartSec=30s
StandardOutput=journal
StandardError=journal
SyslogIdentifier=predict-market-dashboard

[Install]
WantedBy=multi-user.target
